# supervisord-butterfly.conf
[program:butterfly]
command=butterfly.server.py --host="0.0.0.0"
autorestart=true

